# templates/catbot.service.epp
[Unit]
Description=CatBot IRC Bot
After=network.target

[Service]
Type=simple
User=catbot
WorkingDirectory=<%= $catbot::install::home %>
ExecStart=<%= $catbot::install::home %>/env/bin/catbot "<%= $catbot::config::config_file %>"<% if $catbot::config::debug { %> -d<% } %><% if $catbot::config::verbose { %> -v<% } %>
Restart=on-failure

[Install]
WantedBy=multi-user.target
